# conf/supervisor.conf
[program:kappa]
user=root
command=nodejs server.js
directory=/srv/www/kappa
autostart=true
autorestart=true
environment=ENV="live"
